# vulnerable_flask/flaskapp.service
[Unit]
Description=Flask App running on port 6000
After=network.target

[Service]
User=cslab
WorkingDirectory=/
ExecStart=/usr/bin/python3 /home/slab/comps-pentesting/python.py
Restart=always

[Install]
WantedBy=multi-user.target
